# deploy/mesh-rest-bridge.service
[Unit]
Description=Meshtastic BLE-to-JSON bridge (mesh-rest-bridge)
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/share/pac/mesh-gw
ExecStart=/usr/share/pac/mesh-gw/venv/bin/python3 -m cli.main --http-port 8000
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mesh-rest-bridge
TimeoutStopSec=15

# Give bluetoothd time to be fully ready before we try to connect
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
